# scripts/launchd/com.clawd.repo-size-monitor.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.clawd.repo-size-monitor</string>
    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/python3</string>
        <string>/Users/andy/.openclaw/workspace/poc/security/repo_size_monitor.py</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Users/andy/.openclaw/workspace</string>
    <key>StartCalendarInterval</key>
    <dict>
        <key>Weekday</key>
        <integer>1</integer>
        <key>Hour</key>
        <integer>4</integer>
        <key>Minute</key>
        <integer>15</integer>
    </dict>
    <key>StandardOutPath</key>
    <string>/Users/andy/.openclaw/workspace/poc/security/repo_size_stdout.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/andy/.openclaw/workspace/poc/security/repo_size_stderr.log</string>
</dict>
</plist>
