# scripts/launchd/com.clawd.urgent-email.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.clawd.urgent-email</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/andy/.openclaw/workspace/.venv/bin/python3</string>
        <string>/Users/andy/.openclaw/workspace/poc/urgent/poller.py</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Users/andy/.openclaw/workspace</string>
    <key>StartInterval</key>
    <integer>1800</integer>
    <key>StandardOutPath</key>
    <string>/Users/andy/.openclaw/workspace/poc/urgent/poller_stdout.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/andy/.openclaw/workspace/poc/urgent/poller_stderr.log</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PYTHONPATH</key>
        <string>/Users/andy/.openclaw/workspace</string>
    </dict>
</dict>
</plist>
